# worker/crontab
# Luigi worker crontab
*/1 * * * * /worker/run.sh >> /var/log/cron.log 2>&1
# This extra line makes it a valid cron

